# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info
[program:fastapi]
command=/venv/bin/uvicorn main:app --host 0.0.0.0 --port 3001
directory=/app/backend
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/fastapi.stdout.log
stderr_logfile=/var/log/supervisor/fastapi.stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=
    TORCH_CHECKPOINT_PATH="%(ENV_TORCH_CHECKPOINT_PATH)s",
    GPT_MODEL_PATH="%(ENV_GPT_MODEL_PATH)s",
    CORS_ORIGINS="%(ENV_CORS_ORIGINS)s",
    LOG_LEVEL="%(ENV_LOG_LEVEL)s",
    MAX_SESSIONS="%(ENV_MAX_SESSIONS)s",
    SESSION_TTL_HOURS="%(ENV_SESSION_TTL_HOURS)s"
[program:frontend]
command=serve -s /app/frontend/dist -l 8080
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/frontend.stdout.log
stderr_logfile=/var/log/supervisor/frontend.stderr.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
